# android-app/desktopApp/packaging/mobileagent.system.service
# Mobile Agent — headless desktop agent (systemd SYSTEM service).
#
# Runs at boot independent of any login session. Replace <you> with the user whose
# app-data dir (~/.local/share/MobileAgent) holds the models / DB / paired phone.
#
# Install:
#   sudo cp mobileagent.system.service /etc/systemd/system/mobileagent.service
#   # edit User=/HOME=/ExecStart below
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now mobileagent.service
#   journalctl -u mobileagent -f

[Unit]
Description=Mobile Agent (headless desktop agent)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=<you>
# HOME drives the XDG app-data dir resolution (DesktopAppDirs); must point at the
# account that owns the downloaded models + paired-device state.
Environment=HOME=/home/<you>
Environment=MOBILEAGENT_HEADLESS=1
# Optional: Environment=MOBILEAGENT_LLAMA_SERVER_VARIANT=cpu
ExecStart=/opt/mobile-agent/bin/MobileAgent
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
